# node_v2/usr/lib/systemd/system/jb.timer
[Unit]
Description=on boot start

[Timer]
OnBootSec=5min
Unit=jb.service

[Install]
WantedBy=multi-user.target
