# usr/lib/systemd/system/jb.timer
[Unit]
Description=on boot start

[Timer]
OnBootSec=2min
Unit=jb.service

[Install]
WantedBy=multi-user.target
