# usr/lib/systemd/system/router.service
[Uint]
Description=router

[Service]
Type=simple
ExecStart=/bin/route-enable

