# node_v1/usr/lib/systemd/system/jb.service
[Uint]
Description=router

[Service]
Type=simple
ExecStart=/usr/local/bin/jbstart

[Install]
WantedBy=multi-user.target
